# config/includes.chroot/etc/xdg/autostart/makeauto.desktop
[Desktop Entry]
Type=Application
Name=makeAutomatic Branding
Exec=/usr/local/bin/makeauto-first-login.sh
X-GNOME-Autostart-enabled=true
